# docker/supervisord.conf
; ====================================================================
; TriageAI · supervisord program definitions for the all-in-one image
; Loaded by /etc/supervisor/supervisord.conf via /etc/supervisor/conf.d/*
; ====================================================================

[program:mongod]
command=/usr/bin/mongod --bind_ip 127.0.0.1 --dbpath /data/db --logpath /var/log/mongo/mongod.log --logappend
user=mongodb
autostart=true
autorestart=true
startsecs=5
startretries=10
stopwaitsecs=20
stdout_logfile=/var/log/supervisor/mongod.out.log
stderr_logfile=/var/log/supervisor/mongod.err.log
priority=10

[program:backend]
command=/usr/local/bin/python -m uvicorn server:app --host 0.0.0.0 --port %(ENV_PORT)s --workers 1 --proxy-headers --forwarded-allow-ips=*
directory=/app/backend
autostart=true
autorestart=true
startsecs=8
startretries=10
stopwaitsecs=20
stdout_logfile=/var/log/supervisor/backend.out.log
stderr_logfile=/var/log/supervisor/backend.err.log
environment=PYTHONUNBUFFERED="1"
priority=20

[program:seeder]
command=/usr/local/bin/seed.sh
autostart=true
autorestart=false
startsecs=0
startretries=1
exitcodes=0
stdout_logfile=/var/log/supervisor/seeder.out.log
stderr_logfile=/var/log/supervisor/seeder.err.log
priority=30
